# package/mediacenter-osmc/files/DEBIAN/postinst
#!/bin/bash

if [ "$1" = "configure" ]; then
	systemctl enable mediacenter.service >/dev/null 2>&1
fi
