# package/vero3-device-osmc/files/etc/systemd/system/hdr.service
[Unit]
Description=Disables rounding for HDR content

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo round1 > /sys/class/amhdmitx/amhdmitx0/debug'
