# package/vero3-device-osmc/files/etc/systemd/system/hdr.timer
[Unit]
Description = Banding fixup for HDR

[Timer]
OnBootSec=1min
OnUnitActiveSec=1s

[Install]
WantedBy = multi-user.target
