# package/irqbalance-osmc/files/lib/systemd/system/irqbalance.service
[Unit]
Description=irqbalance
After=sysinit.target
 
[Service]
ExecStart=/usr/sbin/irqbalance
Type=forking
 
[Install]
WantedBy=multi-user.target
