# package/rbp2-device-osmc/files/lib/systemd/system/bt-firmware.path
[Unit]
Description=sdio device path monitor

[Path]
PathExists=/dev/serial1
Unit=brcm43xx.service

[Install]
WantedBy=multi-user.target
