# package/a2dp-app-osmc/files/lib/systemd/system/a2dpInfo.service
[Unit]
Description=A2DP Monitor Service

[Install]
WantedBy=multi-user.target
        
[Service]
Type=simple
ExecStart=/usr/bin/a2dpinfo.py 