# package/smb-app-osmc/files/DEBIAN/prerm
#!/bin/bash

if [ "$1" = "remove" ]; then
	systemctl disable samba.service >/dev/null 2>&1
fi

if [ "$1" = "remove" -o "$1" = "upgrade" ]; then
	systemctl stop samba.service
fi
