# package/base-files-osmc/files/etc/systemd/system/emergency.service.d/local.conf
[Service]
ExecStart=
ExecStart=-/bin/sh -c "/usr/bin/setterm --cursor on; /bin/bash; /usr/bin/setterm --cursor off; /bin/systemctl --fail --no-block default"
