# package/transmission-app-osmc/files/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable transmission.service

ischroot
if [ $? == 0 ]; then exit 0; fi

sysctl -p

systemctl start transmission.service
