# package/mediacenter/files/DEBIAN/postinst
#!/bin/bash

ischroot
if [ $? == 0 ]; then exit 0; fi

systemctl enable mediacenter.service
systemctl start mediacenter.service
