# package/vero5-userland-osmc/files/lib/systemd/system/resizemmc.service
[Unit]
ConditionKernelCommandLine=osmc.root_wants_resize=1

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/resizemmc

