# package/network-osmc/files/etc/systemd/system/bluetooth.service.d/local.conf
[Unit]
DefaultDependencies=no
After=dbus.service
Conflicts=shutdown.target
Before=shutdown.target
