# package/vero1-wifi-osmc/files/DEBIAN/postinst
#!/bin/bash

touch /tmp/reboot-needed
systemctl enable vero-wifi.service
