# package/vero1-wifi-osmc/files/DEBIAN/prerm
#!/bin/bash

if [ "$1" = "remove" ]; then
	systemctl disable brcm43xx-firmware.service brcm43xx.service >/dev/null 2>&1
	systemctl stop brcm43xx.service >/dev/null 2>&1
	exit 0
fi

