# package/vero1-wifi-osmc/files/DEBIAN/prerm
#!/bin/bash

systemctl disable vero-wifi.service >/dev/null 2>&1

ischroot
if [ $? == 0 ]; then exit 0; fi

systemctl stop vero-wifi.service
exit 0

