# package/base-files-osmc/files/etc/profile.d/enable-cursor.sh
if tty | egrep -qw "/dev/tty[0-9]+"; then /usr/bin/setterm --cursor on; fi
