# package/network-osmc/files/lib/systemd/system/connman.service
[Unit]
Description=Connection service
Requires=dbus.socket
After=dbus.socket
Wants=network.target
Before=remote-fs-pre.target
Wants=remote-fs-pre.target

[Service]
Type=dbus
BusName=net.connman
Restart=on-failure
ExecStart = /usr/bin/start-network
StandardOutput=null

[Install]
WantedBy=multi-user.target
