# package/irqbalance-osmc/files/DEBIAN/prerm
#!/bin/bash

if [ "$1" = "remove" ]; then
	systemctl disable irqbalance.service >/dev/null 2>&1
	systemctl stop irqbalance.service >/dev/null 2>&1
	exit 0
fi
