# package/diskmount-osmc/files/lib/udev/rules.d/60-cdrom_id.rules
# do not edit this file, it will be overwritten on update

ACTION=="remove", GOTO="cdrom_end"
SUBSYSTEM!="block", GOTO="cdrom_end"
KERNEL!="sr[0-9]*|xvd*", GOTO="cdrom_end"
ENV{DEVTYPE}!="disk", GOTO="cdrom_end"

# unconditionally tag device as CDROM
KERNEL=="sr[0-9]*", ENV{ID_CDROM}="1"

# media eject button pressed - if mounted try to unmount, if not try to eject.
ENV{DISK_EJECT_REQUEST}=="?*", RUN+="/bin/sh -c 'if grep -qw ^$devnode /proc/mounts; then /bin/umount $devnode; else /usr/bin/eject $devnode; fi'", GOTO="cdrom_import"

# import device and media properties and lock tray to
# enable the receiving of media eject button events
IMPORT{program}="cdrom_id --lock-media $devnode", GOTO="cdrom_symlink"

# when the eject button is pressed we still need to import device and media properties
# otherwise udisks will try to unmount disks it shouldn't
LABEL="cdrom_import"
IMPORT{program}="cdrom_id $devnode"

# These rules will create symlinks for the CD/DVD reader, to help old
# programs which are unable to automatically discover the devices.
# The results are undefined for system with multiple CD/DVD devices.

LABEL="cdrom_symlink"
ENV{ID_CDROM}=="?*",        SYMLINK+="cdrom", OPTIONS+="link_priority=-100"
ENV{ID_CDROM_CD_RW}=="?*",  SYMLINK+="cdrw",  OPTIONS+="link_priority=-100"
ENV{ID_CDROM_DVD}=="?*",    SYMLINK+="dvd",   OPTIONS+="link_priority=-100"
ENV{ID_CDROM_DVD_RW}=="?*", SYMLINK+="dvdrw", OPTIONS+="link_priority=-100"

LABEL="cdrom_end"
