# package/base-files-osmc/files/lib/systemd/system/osmc-prune-devel.timer
[Unit]
Description=Timer to prune devel source

[Timer]
OnBootSec=1m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
