# package/base-files-osmc/files/lib/systemd/system/osmc-report-stats.service
[Unit]
Description=Report anonymous OSMC install statistics

[Service]
Type=oneshot
ExecStart=/usr/bin/myosmc-report-installs

[Install]
WantedBy=multi-user.target
