# package/base-files-osmc/files/lib/systemd/system/osmc-prune-devel.service
[Unit]
Description=Service to prune devel source

[Service]
Type=oneshot
ExecStart=/usr/bin/osmc-prune-devel

[Install]
WantedBy=multi-user.target
