# package/atv-remoteclient-osmc/files/DEBIAN/prerm
#!/bin/bash

if [ "$1" = "remove" ]; then
	systemctl disable atvclient.service >/dev/null 2>&1
	systemctl stop atvclient.service >/dev/null 2>&1
	exit 0
fi
