# package/tvheadend-app-osmc/files/DEBIAN/prerm
#!/bin/bash
systemctl disable tvheadend.service >/dev/null 2>&1

ischroot
if [ $? == 0 ]; then exit 0; fi

systemctl stop tvheadend.service
