# package/splash-osmc/files/DEBIAN/prerm
#!/bin/bash
systemctl disable splash.service >/dev/null 2>&1

ischroot
if [ $? == 0 ]; then exit 0; fi

systemctl stop splash.service
