# package/vero5-device-osmc/files/etc/systemd/system/systemd-pstore.service.d/local.conf
[Service]
ExecStartPost=/bin/chown -R osmc:osmc /var/lib/systemd/pstore
