# package/ftr-osmc/files/DEBIAN/postinst
#!/bin/bash
ischroot
if [ $? -eq 0 ]; then systemctl enable ftr.service; fi
