# package/ftr-osmc/files/DEBIAN/postinst
#!/bin/bash
systemctl enable ftr.service
