# package/diskmount-osmc/files/DEBIAN/postinst
#!/bin/bash
echo "postinst 1.0.10"
systemctl enable udisks-glue.service >/dev/null 2>&1

ischroot
if [ $? == 0 ]; then exit 0; fi

systemctl start udisks-glue.service
