# package/a2dp-app-osmc/files/etc/systemd/system/pulseaudio.service
[Unit]
Description=Pulse Audio
        
[Service]
Type=simple
ExecStart=pulseaudio --system --disallow-exit --disallow-module-loading --disable-shm --daemonize