# package/atv-remoteclient-osmc/files/lib/systemd/system/atvclient.service
[Unit]
Description = atvclient
After = remote-fs.target

[Service]
Type = simple
ExecStart = /usr/sbin/atvclient -m
Restart = on-abort

[Install]
WantedBy = multi-user.target
