# package/splash-osmc/files/lib/systemd/system/splash.service
[Unit]
Description = Splash
After = sysinit.target

[Service]
Type=oneshot
ExecStart = /bin/sh -c "export TERM=linux; /usr/bin/setterm --cursor off > /dev/tty1; grep -q quiet /proc/cmdline && /usr/bin/clear > /dev/tty1 && /usr/bin/ply-image; exit 0"

[Install]
WantedBy = basic.target
