# package/aml-vnc-app-osmc/files/DEBIAN/prerm
#!/bin/bash

if [ "$1" = "remove" ]; then
	systemctl disable aml-vnc.service >/dev/null 2>&1
	systemctl stop aml-vnc.service >/dev/null 2>&1
	exit 0
fi
