# package/vero3-device-osmc/files/etc/profile.d/999-device-not-supported.sh
rt() {
echo -e "\033[31m$1\033[0m"
}

rt "Notice: This device is no longer officially supported as of 31/01/2025."
rt "As a result, you will no longer receive officialy supported updates on this device."
rt "You may wish to exercise caution using apt to perform upgrades."
rt "For more information: please see https://osmc.news/QTupI"
