# package/udisks-osmc/files/DEBIAN/prerm
#!/bin/bash

if [ "$1" = "remove" ]; then
	systemctl disable udisks.service >/dev/null 2>&1
	systemctl stop udisks.service >/dev/null 2>&1
	exit 0
fi
