# package/diskmount-osmc/files/lib/systemd/system/udisks-glue.service
[Unit]
Description = mount disks automatically with standby
After = local-fs.target

[Service]
User = osmc
Group = osmc
Type = simple
ExecStartPre = /bin/bash -c "sudo rmdir /media/*/ 2>/dev/null; time udevadm settle; exit 0"
ExecStart = /usr/bin/udisks-glue --foreground

[Install]
WantedBy = multi-user.target
