# package/vero3-wifi-osmc/files/lib/systemd/system/sdio-bt-firmware.path
[Unit]
Description=sdio device path monitor
StopWhenUnneeded=true

[Path]
PathExists=/var/run/sdio_0x02d0:0xa9bf

[Install]
WantedBy=multi-user.target
