# package/vero5-secureosmc-osmc/files/lib/systemd/system/apploader.target
[Unit]
Description=Apploader
StopWhenUnneeded=yes
