# package/vero5-secureosmc-osmc/files/lib/systemd/system/apploader.path
[Unit]
Description=OpTEE path monitoring service
StopWhenUnneeded=true

[Path]
PathExists=/dev/teepriv0

[Install]
WantedBy=multi-user.target
