# package/vero5-secureosmc-osmc/files/lib/systemd/system/tee.mount
[Unit]
Description=OPTEE secure data mount

[Mount]
What=/dev/tee
Where=/tee
Type=ext4
Options=rw

[Install]
WantedBy=multi-user.target
