# package/eventlircd-osmc/files/DEBIAN/postinst
#!/bin/bash

systemctl enable eventlircd.service
