# package/eventlircd-osmc/files/DEBIAN/postinst
#!/bin/bash

systemctl enable eventlircd.service

ischroot
if [ $? == 0 ]; then exit 0; fi

systemctl start eventlircd.service

