# package/eventlircd-osmc/files/DEBIAN/prerm
#!/bin/bash

if [ "$1" = "remove" ]; then
	systemctl disable eventlircd.service >/dev/null 2>&1
	systemctl stop eventlircd.service >/dev/null 2>&1
	exit 0
fi
