# package/mediacenter-next-osmc/files/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable mediacenter.service
