# package/mediacenter-next-osmc/files/DEBIAN/postinst
#!/bin/bash

systemctl enable mediacenter.service
